# files/docker.service.d/selinux.conf
[Service]
Environment=DOCKER_SELINUX=""
