# files/docker.service.d/selinux.conf
[Service]
# Vi bruker foreløpig ikke SELinux
Environment=DOCKER_SELINUX=""
