# templates/flannel.service.j2
{% set etcd_urls = [] -%}
{% for host in groups['etcd'] -%}
  {% set _ = etcd_urls.append("https://%s:2379" % (host)) -%}
{% endfor -%}
[Unit]
Description=flannel
Documentation=https://coreos.com/flannel
Before=docker.service

[Service]
Environment=GODEBUG=netdns=go
Environment=FLANNEL_IMAGE_TAG=v{{ flannel_version }}
Environment=HTTPS_PROXY={{ nais_https_proxy }}
Environment=NO_PROXY={{ nais_no_proxy }}
ExecStart=/usr/lib/flatcar/flannel-wrapper \
  --etcd-endpoints={{ etcd_urls | join(",") }} \
  --etcd-prefix=/nais/network \
  --etcd-cafile=/etc/ssl/certs/ca.pem \
  --etcd-certfile=/etc/ssl/certs/etcd-client.pem \
  --etcd-keyfile=/etc/ssl/certs/etcd-client-key.pem
Restart=always
StartLimitInterval=0
RestartSec=10

[Install]
WantedBy=multi-user.target
