# files/etcd/etcd-prometheus-push.timer
[Unit]
Description=Push metrics to prometheus every minute

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=multi-user.target
