# files/etcd/etcd-backup.timer
[Unit]
Description=Run etcd-backup service every night

[Timer]
OnCalendar=*-*-* 03:00:00
