# files/etcd/etcd-backup.timer
[Unit]
Description=Run etcd-backup service every night

[Timer]
OnCalendar=*-*-* 22:00:00

[Install]
WantedBy=multi-user.target

