# templates/flanneld.service.d/etcd.conf
{% set etcd_urls = [] -%}
{% for host in groups['etcd'] -%}
  {% set _ = etcd_urls.append("https://%s:2379" % (host)) -%}
{% endfor -%}

[Service]
Environment=GODEBUG=netdns=go
Environment="FLANNELD_ETCD_ENDPOINTS={{ etcd_urls | join(",") }}"
Environment="FLANNELD_ETCD_PREFIX=/nais/network"
Environment="FLANNELD_ETCD_CAFILE=/etc/ssl/certs/ca.pem"
Environment="FLANNELD_ETCD_CERTFILE=/etc/ssl/certs/etcd-client.pem"
Environment="FLANNELD_ETCD_KEYFILE=/etc/ssl/certs/etcd-client-key.pem"
