# templates/flanneld.service.d/etcd.conf
{% set etcd_urls = [] -%}
{% for host in groups['etcd'] -%}
  {% set _ = etcd_urls.append("https://%s:2379" % (host)) -%}
{% endfor -%}

[Service]
Environment="FLANNEL_OPTS=--ip-masq=true --etcd-endpoints={{ etcd_urls | join(",") }} --etcd-prefix=/nais/network --etcd-cafile=/etc/ssl/certs/ca.pem --etcd-certfile=/etc/ssl/certs/etcd-client.pem --etcd-keyfile=/etc/ssl/certs/etcd-client-key.pem"