# templates/flanneld.service.d/etcd.conf
{% set etcd_urls = [] -%}
{% for host in groups['etcd'] -%}
  {% set _ = etcd_urls.append("https://%s:2379" % (host)) -%}
{% endfor -%}

[Service]
Environment=GODEBUG=netdns=go
Environment="FLANNELD_KUBE_API_URL=https://b27apvl00534.preprod.local:6443"
Environment="FLANNELD_KUBE_SUBNET_MGR=true"
Environment="FLANNELD_IP_MASQ=true"
Environment="FLANNELD_IFACE=ens192"
Environment="NODE_NAME={{ inventory_hostname }}"
# flannel-wrapper script only support mounting ETCD_SSL_DIR, use this to mount dir with kube-config
Environment="ETCD_SSL_DIR=/etc/kubernetes"
Environment="FLANNELD_KUBECONFIG_FILE=/etc/kubernetes/kubeconfigs/kubelet.conf"
Environment="FLANNELD_NET_CONFIG_PATH=/etc/kubernetes/flannel/net-conf.json"
